# kanjuro/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/storage/logs/supervisor.log

[program:php]
command = php-fpm
